# obmc-webserver.socket
[Unit]
Description=BMC Web Server socket

[Socket]
ListenStream=443
ReusePort=true

[Install]
WantedBy=sockets.target
